# .config/systemd/user/clipmenud.service.d/override.conf
[Service]
Environment=CM_DIR=%h/.cache/
Environment=XAUTHORITY=/run/user/1002/Xauthority
